# pkgs/alpine/crontab
@reboot syncthing serve --no-default-folder
